# social-automation/launchd/com.socialautomation.daily-wp-draft.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.socialautomation.daily-wp-draft</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/gilcohen/.pyenv/versions/3.11.9/bin/python3</string>
        <string>scripts/daily_wp_draft.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/gilcohen/Projects/dogfoodandfun/social-automation</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/gilcohen/Projects/dogfoodandfun/social-automation/logs/cron_daily_wp_draft.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gilcohen/Projects/dogfoodandfun/social-automation/logs/cron_daily_wp_draft.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>BRAND_DIR</key>
        <string>/Users/gilcohen/Projects/dogfoodandfun/dogfoodandfun</string>
        <key>PATH</key>
        <string>/Users/gilcohen/.pyenv/versions/3.11.9/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>HOME</key>
        <string>/Users/gilcohen</string>
    </dict>
</dict>
</plist>
