# social-automation/launchd/com.socialautomation.ig-own-comments.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.socialautomation.ig-own-comments</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/gilcohen/.pyenv/versions/3.11.9/bin/python3</string>
        <string>scripts/ig_own_comments.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/gilcohen/Projects/dogfoodandfun/social-automation</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/gilcohen/Projects/dogfoodandfun/social-automation/logs/cron_ig_own_comments.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gilcohen/Projects/dogfoodandfun/social-automation/logs/cron_ig_own_comments.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>BRAND_DIR</key>
        <string>/Users/gilcohen/Projects/dogfoodandfun/dogfoodandfun</string>
        <key>PATH</key>
        <string>/Users/gilcohen/.pyenv/versions/3.11.9/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>HOME</key>
        <string>/Users/gilcohen</string>
    </dict>
</dict>
</plist>
